# root/etc/systemd/system/rspamd.service.d/clamd.conf
#start clamd by rspamd.service
[Unit]
Wants=clamd@rspamd.service

[Service]

